# service/snapper-boot-entries.path
[Unit]
Description=Watch for snapper snapshot changes

[Path]
PathChanged=/.snapshots
Unit=snapper-boot-entries.service
TriggerLimitIntervalSec=5min

[Install]
WantedBy=multi-user.target
